# xdoc/check_service.cron
*/1 * * * * ./check_service.sh zd 8849 >> check_service.log